# debian/jetty9.timer
[Unit]
Description=Weekly trigger for jetty9 logging. 

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
